# mullvad-speed-guard/launchagents/com.story.mullvad-speed-guard.panel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.story.mullvad-speed-guard.panel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__HOME__/Library/Application Support/MullvadSpeedGuard/guard_panel_server.py</string>
    <string>--port</string>
    <string>18790</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__HOME__/Library/Application Support/MullvadSpeedGuard</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>StandardOutPath</key>
  <string>__HOME__/Library/Application Support/MullvadSpeedGuard/results/panel_launchagent.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Application Support/MullvadSpeedGuard/results/panel_launchagent.log</string>
</dict>
</plist>
